# systemd/fpren-multi-zone-streamer.service
[Unit]
Description=FPREN Multi-Zone Streamer — one FFmpeg → Icecast stream per zone
After=network.target icecast2.service mongod.service zone-alert-tts.service
Wants=icecast2.service mongod.service zone-alert-tts.service

[Service]
Type=simple
User=ufuser
# WorkingDirectory must be the repo root (not weather_station/) so that
# "weather_station.*" package imports resolve without a pip install.
WorkingDirectory=/home/ufuser/Fpren-main
ExecStart=/home/ufuser/Fpren-main/venv/bin/python3 -m weather_station.services.multi_zone_streamer
Restart=on-failure
RestartSec=30
TimeoutStopSec=60
EnvironmentFile=/home/ufuser/Fpren-main/weather_station/.env
StandardOutput=append:/home/ufuser/Fpren-main/logs/multi_zone_streamer.log
StandardError=append:/home/ufuser/Fpren-main/logs/multi_zone_streamer.log

[Install]
WantedBy=multi-user.target
