# systemd/fpren-census-refresh.timer
[Unit]
Description=Refresh FPREN Florida Census data monthly

[Timer]
Unit=beacon-census-fetcher.service
OnCalendar=monthly
AccuracySec=1d
Persistent=true

[Install]
WantedBy=timers.target
