# systemd/zone-alert-tts.service
[Unit]
Description=FPREN Zone Alert TTS — transcodes NWS alerts and traffic to MP3 per zone
After=network.target mongod.service beacon-weather-service.service
Wants=mongod.service

[Service]
Type=simple
User=ufuser
WorkingDirectory=/home/ufuser/Fpren-main
ExecStart=/home/ufuser/Fpren-main/venv/bin/python3 -m weather_station.services.zone_alert_tts
Restart=on-failure
RestartSec=15
TimeoutStopSec=150
EnvironmentFile=/home/ufuser/Fpren-main/weather_station/.env

[Install]
WantedBy=multi-user.target
